# files/minesweep-systemd.timer
[Unit]
Description="Run minesweep-systemd.sh every 30s"

[Timer]
OnActiveSec=30
Persistent=true

[Install]
WantedBy=timers.target
