# files/systemd-minesweep.timer
[Unit]
Description="Run systemd-minesweep.sh every 30s"

[Timer]
OnActiveSec=30s
Persistent=true

[Install]
WantedBy=timers.target
